# scripts/dnsmasq-web.service
[Unit]
Description=dnsmasq-web — web management console for dnsmasq
Documentation=file:///opt/dnsmasq-web/README.md
Wants=network-online.target
After=network-online.target dnsmasq.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/dnsmasq-web
ExecStart=/opt/dnsmasq-web/dnsmasq-web

# ── Configuration ────────────────────────────────────────────────────
# HOST=127.0.0.1 keeps the console local-only. Set HOST= (empty) to listen
# on all interfaces — enable AUTH_PASSWORD below (or put an authenticating
# reverse proxy in front) before you do.
Environment=HOST=127.0.0.1
Environment=PORT=8053
Environment=DNSMASQ_CONF=/etc/dnsmasq.conf
Environment=DNSMASQ_LEASES=/var/lib/dnsmasq/dnsmasq.leases
Environment=BACKUP_DIR=/var/backups/dnsmasq-web

# Optional login. Set AUTH_PASSWORD to require a sign-in for the console;
# set API_TOKEN so non-browser clients (the MCP server, curl) can use
# `Authorization: Bearer $API_TOKEN`. Enable BOTH if an AI agent should
# keep working after auth is turned on.
#
# NOTE: unit files are world-readable (644). On a single-admin box that's
# fine; on a shared machine put the secrets in a root-only file instead:
#   printf 'AUTH_PASSWORD=…\nAPI_TOKEN=…\n' > /etc/dnsmasq-web.env
#   chmod 600 /etc/dnsmasq-web.env
# and uncomment EnvironmentFile below rather than the Environment lines.
#Environment=AUTH_PASSWORD=change-me
#Environment=API_TOKEN=generate-a-long-random-string
#EnvironmentFile=/etc/dnsmasq-web.env

Restart=on-failure
RestartSec=3

# ── Sandboxing ───────────────────────────────────────────────────────
# Needs root for /etc/dnsmasq.conf writes + systemctl, but everything
# below is free hardening that doesn't get in the way of either.
PrivateTmp=yes
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ProtectClock=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
